# src/etc/cron.d/airpuff-html.orig
#*/5 * * * * root /usr/local/bin/airpuff/airpuff.sh
#0 00,12,13,14,15,16,23 * * 1-5 root /usr/local/bin/airpuff/airpuff-cron.sh

# Quick commute Visi chart:
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd-commute-vis.sh
# Get the wx
*/5 * * * * root /usr/local/bin/airpuff/airpuff-custom.sh "Work Commute" "KEDU KVCB KSUU KCCR KLVK KHWD KOAK KRHV KSJC KSQL KNUQ KPAO"
*/5 * * * * root /usr/local/bin/airpuff/airpuff-custom.sh "Sac Valley" "KEDU KSAC KWLW KCIC KOVE KRBL KRDD KBAB KGOO"
*/5 * * * * root /usr/local/bin/airpuff/airpuff-custom.sh "National Airports" "KSMF KSFO KLAX KSEA KPDX KDEN KIAH KORD KMSP KJFK KLGA KIAD KCLT KATL KMSY KBOS"
# Store the data
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd.sh "KEDU KVCB KSUU KCCR KLVK KHWD KOAK KRHV KSJC KSQL KNUQ KPAO"
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd.sh "KSAC KWLW KCIC KOVE KRBL KRDD KBAB KGOO"
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd.sh "KSMF KSFO KLAX KSEA KPDX KDEN KIAH KORD KMSP KJFK KLGA KIAD KCLT KATL KMSY KBOS"
# Make the charts
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd-graph.sh "KEDU KVCB KSUU KCCR KLVK KHWD KOAK KRHV KSJC KSQL KNUQ KPAO"
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd-graph.sh "KSAC KWLW KCIC KOVE KRBL KRDD KBAB KGOO"
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd-graph.sh "KSMF KSFO KLAX KSEA KPDX KDEN KIAH KORD KMSP KJFK KLGA KIAD KCLT KATL KMSY KBOS"
# Make the web pages
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd-web-create.sh "KEDU KVCB KSUU KCCR KLVK KHWD KOAK KRHV KSJC KSQL KNUQ KPAO"
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd-web-create.sh "KSAC KWLW KCIC KOVE KRBL KRDD KBAB KGOO"
*/5 * * * * root /usr/local/bin/airpuff/airpuff-rrd-web-create.sh "KSMF KSFO KLAX KSEA KPDX KDEN KIAH KORD KMSP KJFK KLGA KIAD KCLT KATL KMSY KBOS"
# Cleanup
6,11,16,21,26,31,36,41,46,51,56 * * * * root /usr/local/bin/airpuff/airpuff-rrd-web-create-cleanup.sh

