# src/etc/cron.d/clean-root-email
*/30 * * * * root /bin/cat /dev/null > /var/spool/mail/root

