# ansible/roles/airpuff_app/templates/airpuff-fetcher.timer.j2
[Unit]
Description=Run AirPuff Fetcher every 5 minutes
Requires=airpuff-fetcher.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
