# ansible/roles/airpuff_app/templates/airpuff-fetcher.service.j2
[Unit]
Description=AirPuff Fetcher
After=network.target postgresql.service redis.service airpuff-web.service
Wants=postgresql.service redis.service airpuff-web.service

[Service]
Type=exec
User={{ airpuff_app_user }}
Group={{ airpuff_app_group }}
WorkingDirectory={{ airpuff_app_home }}/backend
Environment=PATH={{ virtualenv_path }}/bin
Environment=PYTHONPATH={{ airpuff_app_home }}/backend
ExecStart={{ virtualenv_path }}/bin/python -m app.services.weather_fetcher
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=airpuff-fetcher

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ airpuff_app_home }} {{ airpuff_app_logs }} {{ airpuff_app_data }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target
